# phpsessionclean/phpsessionclean2.timer
[Unit]
Description=Clean PHP session files every 30 mins - 2

[Timer]
OnCalendar=*-*-* *:09,39:00
Persistent=true

[Install]
WantedBy=timers.target
