# phpsessionclean/phpsessionclean2.service
[Unit]
Description=Clean php session files

[Service]
Type=oneshot
ExecStart=/usr/lib/php/sessionclean2
ProtectHome=true
ProtectSystem=true
PrivateTmp=true
